# base-files/etc/uci-defaults/100_uci_custom
uci set network.lan.dns="127.0.0.1 8.8.8.8"
uci set network.wan.peerdns="0"
uci set network.wan.dns="127.0.0.1 8.8.8.8"
uci set network.wan6.peerdns="0"
uci set network.wan6.dns="::1 2001:4860:4860::8888"
uci commit network
uci set network.lan.network="docker0"
uci commit network